# tattood.service
[Unit]
Description=uWSGI instance to serve tattood
After=network.target

[Service]
User=root
Group=root
WorkingDirectory=/root/tattood/
Environment="PATH=/root/tattood/venv/bin"
ExecStart=/root/tattood/venv/bin/uwsgi --ini tattood.ini

[Install]
WantedBy=multi-user.target
